# deploy/host-watcher/msm-selfupdate.service
[Unit]
Description=MSM Host-Watcher: führt angefordertes Selbst-Update aus

[Service]
Type=oneshot
ExecStart=/usr/local/bin/msm-selfupdate-watcher.sh
