# deploy/host-watcher/msm-upgrade.service
[Unit]
Description=MSM Host-Watcher: führt angeforderten MC-Versionssprung aus

[Service]
Type=oneshot
ExecStart=/usr/local/bin/msm-upgrade-watcher.sh
